# atari-scripts/atari-high-score-rotate.crontab
0 0 1 * * thomc /usr/local/bin/atari-high-score-rotate.sh

